# IPTV_v5.4_MAX_AGGRESSION/vps/network/brazil-route-persist.path
[Unit]
Description=Watch wg-surfshark-br interface state changes
[Path]
PathChanged=/etc/wireguard/wg-surfshark-br.conf
Unit=brazil-route-persist.service
[Install]
WantedBy=multi-user.target
