# IPTV_v5.4_MAX_AGGRESSION/vps/network/wg-health-fast.service
[Unit]
Description=WG Health Fast (sub-second failover detector)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/netshield/scripts/wg-health-fast.sh
StandardOutput=null
StandardError=journal
