# IPTV_v5.4_MAX_AGGRESSION/net-shield/snapshots/2026-04-26_funcional_post_rollback/etc/systemd/system/netshield-autopilot.timer
[Unit]
Description=Run NET SHIELD Autopilot every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Unit=netshield-autopilot.service

[Install]
WantedBy=timers.target
